# reminder-daemon.service
[Unit]
Description=Benten Schedule Reminder Daemon
After=network.target

[Service]
Type=simple
User=connoise
WorkingDirectory=/home/connoise/Schedule
ExecStart=/usr/bin/python3 /home/connoise/Schedule/reminder_daemon.py
Restart=always
RestartSec=10
Environment=TELEGRAM_BOT_TOKEN=8313726439:AAHDDrgV-nGqvhJSam5eNuBu9pLi3ScrdBQ
Environment=TELEGRAM_CHAT_ID=-5218106664
Environment=REMINDERS_FILE=/home/connoise/Schedule/reminders.json

[Install]
WantedBy=multi-user.target
